# 03_playbooks/Example_1/templates/flask_supervisor.conf
[program:app]
command={{ webapps_dir }}/{{ app_name }}/app/{{ app_name }}
user=vagrant
autostart=true
autorestart=true
stdout_logfile=/var/log/app.log
stderr_logfile=/var/log/app.err.log
